# fsback.service
[Unit]
Description=fsback service

[Service]
Type=oneshot
ExecStart=/usr/bin/fsback --sysdaemon
